# profile.d/bash_aliases.sh
#@IgnoreInspection BashAddShebang

# Directories #
alias ll='ls -alF'
alias la='ls -A'
alias l='ls -CF1'

# Searching #
alias g='grep -sIrEn'
alias qf='find . -name'
alias qfi='find . -iname'

alias fw='date +FW%V.%u'

#alias hp='$(history | peco | cut -d " " -f 4-)'

alias t='tree'
